# deployment/storm-scout-dev.service
[Unit]
Description=Storm Scout DEV/QC Server
Documentation=https://github.com/murphy-platforms/storm-scout-usps
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=your-ssh-user
WorkingDirectory=/path/to/storm-scout/backend
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=storm-scout-dev

# Environment
Environment=NODE_ENV=development

# Limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
